# raid_check.service
[Unit]
Description=Backup Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/raid_check
